# crontab
# START CRON JOB
* * * * * /app/venv/bin/python3 /app/subscribe.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
